# .infrastructure/crontab
0 0 * * * docker exec -t y2b-php-fpm bash -c 'bin/console landing:sitemap:generate -vv'